# systemd/vault-search-reindex.service
[Unit]
Description=Reindex Obsidian vaults for semantic search
Documentation=file:///home/rtm/wrk/projects/obsidian-mcp/README.md
# Ollama is a system service, so it cannot be a hard dependency of a user unit.
# reindex.py exits non-zero with a clear message when Ollama is unreachable and
# the timer simply retries on the next tick.

[Service]
Type=oneshot
WorkingDirectory=/home/rtm/wrk/projects/obsidian-mcp
ExecStart=/home/rtm/wrk/projects/obsidian-mcp/.venv/bin/python reindex.py -q
# First full index embeds ~23k chunks on CPU and can run well over an hour.
TimeoutStartSec=4h
Nice=10
IOSchedulingClass=idle
